# deploy/systemd/atlas-caddy.service
[Unit]
Description=Atlas Caddy Reverse Proxy (HTTPS + TLS)
After=network.target
Wants=atlas-oauth2-proxy.service atlas-telemetry.service

[Service]
Type=simple
ExecStart=/usr/bin/caddy run --config /home/claude/Caddyfile
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
