# deploy/systemd/atlas-superego.service
[Unit]
Description=Atlas Superego (Gemma 4 31B on GPU 0)
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=claude
Environment=CUDA_VISIBLE_DEVICES=0
ExecStart=/home/claude/llama.cpp/build/bin/llama-server \
    --model /home/claude/models/gemma-4-31B-it-Q5_K_M/gemma-4-31B-it-Q5_K_M.gguf \
    --host 0.0.0.0 --port 8080 --ctx-size 8192 --threads 12
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=atlas.target
